# scripts/nocobase-briefing-worker.service
[Unit]
Description=NocoBase Briefing Worker
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/bin/node /opt/noco-base/worker/briefing-worker.js
Restart=always
RestartSec=10
Environment=NOCOBASE_BASE=https://localhost:668
Environment=NOCOBASE_EMAIL=voadge@voadge.cn
Environment=NOCOBASE_PASSWORD=admin123
WorkingDirectory=/opt/noco-base/worker
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
